# supervisor.ini
[supervisord]
nodaemon=true

[program:php-download]
process_name=%(program_name)s_%(process_num)02d
command = php app/converter.php
user = root
autostart = true
autorestart = true
numprocs=3

[program:php-upload]
process_name=%(program_name)s_%(process_num)02d
command = php app/upload.php
user = root
autostart = true
autorestart = true
numprocs=3
